# conda/config/entityd.supervisord.conf
{# Supervisor configuration template for entityd. -#}

[program:entityd]
{%- if 'EDD_X_MODELD' in environ %}
command = %PREFIX%/bin/entityd --dest {{environ['EDD_X_MODELD']}}
{% else %}
command = %PREFIX%/bin/entityd
{% endif -%}
redirect_stderr = true
