# conda/entityd.supervisord.conf
# supervisord configuration for entityd

[program:entityd]
command = %PREFIX%/bin/entityd